# deploy/com.quasarryan.mlxaudio.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quasarryan.mlxaudio.api</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ryan/works/services/mlx-audio-bridge/.venv/bin/mlx-audio-bridge-server</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ryan/works/services/mlx-audio-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>API_KEY</key>
    <string>change-me</string>
    <key>QWEN_TTS_MODEL</key>
    <string>mlx-community/Qwen3-TTS-12Hz-0.6B-Base-bf16</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ryan/works/services/mlx-audio-bridge/run/mlx-audio-bridge.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ryan/works/services/mlx-audio-bridge/run/mlx-audio-bridge.stderr.log</string>
</dict>
</plist>
