# deploy/com.quasarryan.mlxaudio.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quasarryan.mlxaudio.api</string>

  <!-- Default: run as root in the privileged system domain. -->
  <key>UserName</key>
  <string>root</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/mlx-audio-bridge/.venv/bin/mlx-audio-bridge-server</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/opt/mlx-audio-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>API_KEY</key>
    <string>change-me</string>
    <key>BIND_ADDRESS</key>
    <string>127.0.0.1</string>
    <key>PORT</key>
    <string>8008</string>
    <key>QWEN_MODEL_DIR</key>
    <string>/opt/mlx-audio-bridge/models</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/opt/mlx-audio-bridge/run/mlx-audio-bridge.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/opt/mlx-audio-bridge/run/mlx-audio-bridge.stderr.log</string>
</dict>
</plist>
